# bootstrap/supervisord.conf
[program:lfs_lab_cert_tracker]
command=/usr/local/bin/gunicorn -b 0.0.0.0:8080 lfs_lab_cert_tracker.wsgi
directory=/vagrant
autostart=true
autorestart=true
stderr_logfile=/var/log/lfs_lab_cert_tracker/lfs_lab_cert_tracker.err.log
stdout_logfile=/var/log/lfs_lab_cert_tracker/lfs_lab_cert_tracker.out.log
